# res/systemd/rss.timer
[Unit]
Description=Hourly refresh of singlerss

[Timer]
OnCalendar=hourly
 
[Install]
WantedBy=timers.target
